# roles/deploy_dss_benchmark/templates/simplehttp.service.j2
[Unit]
Description=python SimpleHTTPServer daemon
Documentation=man:SimpleHTTPServer(1)

[Service]
Type=simple
WorkingDirectory={{ dss_log_dir }}
ExecStart=/usr/bin/python -m SimpleHTTPServer {{ deploy_dss_benchmark_httpd_port }}

[Install]
WantedBy=multi-user.target
